# XOSAUTOCONFIG/CORE.localhost/xosautoconfig/etc/profile.d/proxy.csh

setenv http_proxy http://paradonf.irisa.fr:3128
setenv https_proxy http://paradonf.irisa.fr:3128
setenv ftp_proxy http://paradonf.irisa.fr:3128
