# XOSAUTOCONFIG/NODE_VBOX1/xosautoconfig/etc/profile.d/proxy.sh

http_proxy=http://paradonf.irisa.fr:3128
https_proxy=http://paradonf.irisa.fr:3128
ftp_proxy=http://paradonf.irisa.fr:3128
